# crontab
* * * * * "/usr/local/bin/php" "/app/api/artisan" schedule:run >>/app/ezwork_api.log
 